# mklinks.service
[Unit]

[Service]
ExecStart=/usr/local/sbin/mklinks

[Install]
WantedBy=default.target
